# systemd/whim-server.service
[Unit]
Description=whim-proxy webhook relay server
After=network.target

[Service]
ExecStart=/usr/local/bin/whim-server --addr :9000
Restart=on-failure
RestartSec=5s

User=whim
Group=whim

# Hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
CapabilityBoundingSet=

[Install]
WantedBy=multi-user.target
